# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

; ── FastAPI (port 7860 — HF Spaces default) ──────────────────
[program:fastapi]
command=uvicorn app.api:app --host 0.0.0.0 --port 7860 --workers 1
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/tmp/fastapi.log
stderr_logfile=/tmp/fastapi_err.log
stdout_logfile_maxbytes=10485760
stderr_logfile_maxbytes=10485760

; ── Streamlit Dashboard (port 8501) ──────────────────────────
[program:streamlit]
command=streamlit run app_dashboard.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/tmp/streamlit.log
stderr_logfile=/tmp/streamlit_err.log
stdout_logfile_maxbytes=10485760
stderr_logfile_maxbytes=10485760