# supervisord.local.conf
; ══════════════════════════════════════════════════════════════
; supervisord.local.conf  —  LOCAL DEVELOPMENT ONLY
;
; Manages both Streamlit apps inside the local Docker container.
; Used by Dockerfile.local. NOT used on HuggingFace.
;
; Fraud Detection  →  port 8501
; Universal        →  port 8502
; ══════════════════════════════════════════════════════════════

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

; ── Fraud Detection App (port 8501) ──────────────────────────────
[program:fraud]
command=streamlit run app_fraud_local.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/tmp/fraud_stdout.log
stderr_logfile=/tmp/fraud_stderr.log
stdout_logfile_maxbytes=10485760
stderr_logfile_maxbytes=10485760
environment=PYTHONUNBUFFERED="1"

; ── Universal Trainer App (port 8502) ────────────────────────────
[program:universal]
command=streamlit run app_universal.py
    --server.port=8502
    --server.address=0.0.0.0
    --server.headless=true
    --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/tmp/universal_stdout.log
stderr_logfile=/tmp/universal_stderr.log
stdout_logfile_maxbytes=10485760
stderr_logfile_maxbytes=10485760
environment=PYTHONUNBUFFERED="1"