# supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

; ── FastAPI (main API — HF Spaces default port 7860) ─────────
[program:fastapi]
command=uvicorn app.api:app --host 0.0.0.0 --port 7860 --workers 1
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/fastapi.log
stderr_logfile=/app/logs/fastapi_err.log
stdout_logfile_maxbytes=10485760
stderr_logfile_maxbytes=10485760

; ── Streamlit Dashboard (port 8501) ──────────────────────────
[program:streamlit]
command=streamlit run app/dashboard.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/streamlit.log
stderr_logfile=/app/logs/streamlit_err.log
stdout_logfile_maxbytes=10485760
stderr_logfile_maxbytes=10485760

; ── MLflow UI (port 5000) ─────────────────────────────────────
[program:mlflow]
command=mlflow ui --host 0.0.0.0 --port 5000 --backend-store-uri sqlite:///app/models/mlflow.db
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/app/logs/mlflow.log
stderr_logfile=/app/logs/mlflow_err.log
stdout_logfile_maxbytes=10485760
stderr_logfile_maxbytes=10485760